# installer/debian/altcodes4linux/DEBIAN/postinst
#!/bin/bash

# exit on error
set -e

# source debconf library.
#. /usr/share/debconf/confmodule

systemctl daemon-reload
systemctl enable altcodes4linux
systemctl start altcodes4linux
