# installer/debian/altcodes4linux/DEBIAN/prerm
#!/bin/bash

# exit on error
set -e

# source debconf library.
#. /usr/share/debconf/confmodule

systemctl stop altcodes4linux
systemctl disable altcodes4linux
